# deploy/crontab/crontab
13 20 * * * ansible-playbook main.yaml -i inventory -l dbservers -b
